# varnish/varnishd.conf
[program:varnishd]
command=/usr/sbin/varnishd -F -f /etc/varnish/main.vcl
priority=0
autostart=true
autorestart=true

