# varnish/varnishncsa.conf
[program:varnishncsa]
command=/usr/bin/varnishncsa
priority=1
autostart=true
autorestart=true

